# systemd/.config/systemd/user/rsync-backup-podman.timer
[Unit]
Description=Run daily backups

[Timer]
OnCalendar=*-*-* 2:00
Persistent=true

[Install]
WantedBy=timers.target
