# systemd/.config/systemd/user/rsync-backup-podman.service
[Unit]
Description=Back up podman volumes

[Service]
Type=oneshot
ExecStart=%h/bin/rsync-backup-podman.sh
