# podman/.config/containers/systemd/cloudflare.container
[Unit]
Description=Cloudflare
Wants=network-online.target
After=network-online.target
After=local-fs.target

[Container]
ContainerName=cloudflare
Image=docker.io/cloudflare/cloudflared:latest
AutoUpdate=registry
UserNS=keep-id
GroupAdd=keep-groups
User=1000:1000
Network=host
Secret=cloudflare_tunnel_token,type=env,target=TUNNEL_TOKEN
Exec=tunnel --no-autoupdate run

[Service]
Restart=on-failure
TimeoutStartSec=900

[Install]
WantedBy=default.target
