# podman/.config/containers/systemd/emby.container
[Unit]
Description=EmbyServer
After=local-fs.target network.target

[Container]
ContainerName=embyserver
Image=lscr.io/linuxserver/emby:latest
#Image=emby/embyserver:latest
AutoUpdate=registry
Network=host
Environment=TZ=Europe/Stockholm
Volume=embyserver_data:/config
Volume=/mnt/media:/media
Volume=/mnt/backup/emby:/backup
AddDevice=/dev/dri

[Service]
Restart=on-failure
TimeoutStartSec=900

[Install]
#WantedBy=default.target
