# podman/.config/containers/systemd/lidarr.container
[Unit]
Description=Lidarr container
Wants=network-online.target
After=network-online.target local-fs.target

[Container]
ContainerName=lidarr
Image=lscr.io/linuxserver/lidarr:latest
AutoUpdate=registry
UserNS=keep-id
GroupAdd=keep-groups
Pod=vpn_stack.pod
Volume=lidarr_data:/config
Volume=/mnt/media/Musik:/music
Volume=/mnt/externa/Nedladdningar:/downloads
Network=gluetun.container

[Service]
Restart=on-failure
TimeoutStartSec=900

[Install]
WantedBy=default.target
